# .config/autostart/Handy.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Handy
Comment=Handystartup script
Exec=/opt/Handy.AppImage 
StartupNotify=false
Terminal=false